# flaskapp_supervisor.conf
[program:app-uwsgi]
command = /usr/local/bin/uwsgi --ini /data/rutaaccesible_uwsgi.ini

[program:nginx-app]
command = /usr/sbin/nginx